# payload/Library/LaunchAgents/com.erikng.installapplications.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.erikng.installapplications</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Library/installapplications/Python.framework/Versions/3.8/bin/python3</string>
		<string>/Library/installapplications/installapplications.py</string>
		<string>--userscript</string>
	</array>
	<key>KeepAlive</key>
  <dict>
     <key>PathState</key>
     <dict>
       <key>/var/tmp/installapplications/.userscript</key>
       <true/>
     </dict>
  </dict>
	<key>OnDemand</key>
	<true/>
	<key>StandardOutPath</key>
	<string>/var/log/installapplications/installapplications.user.log</string>
	<key>StandardErrorPath</key>
	<string>/var/log/installapplications/installapplications.user.log</string>
</dict>
</plist>
