# payload/Library/LaunchDaemons/com.erikng.installapplications.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.erikng.installapplications</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Library/installapplications/Python.framework/Versions/Current/bin/python3</string>
		<string>/Library/installapplications/installapplications.py</string>
		<string>--jsonurl</string>
		<string>https://domain.tld</string>
		<!-- <string>--dry-run</string> -->
		<!-- <string>--follow-redirects</string> -->
		<!-- <string>--headers</string> -->
		<!-- <<string>Basic ZXhhbXBsZV91c2VyOmV4YW1wbGVfcGFzc3dvcmQ=</string> -->
		<!-- <string>--iapath</string> -->
		<!-- <string>/Library/installapplications/installapplications</string> -->
		<!-- <string>--laidentifier</string> -->
		<!-- <string>com.erikng.installapplications</string> -->
		<!-- <string>--ldidentifier</string> -->
		<!-- <string>com.erikng.installapplications</string> -->
		<!-- <string>--reboot</string> -->
		<!-- <string>--skip-validation</string> -->
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardOutPath</key>
	<string>/var/log/installapplications/installapplications.log</string>
	<key>StandardErrorPath</key>
	<string>/var/log/installapplications/installapplications.log</string>
</dict>
</plist>
